# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:tellstickcluster]
command=java -Dip=%(ENV_ip)s -jar tellstick-hazelcast-cluster.jar

[program:telldusd]
command=telldusd --nodaemon
